# qr3_95.conf
[program:qr3_95]
command=python3 /home/curator3/qr3_95/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_95/test.err.log
stdout_logfile=/home/curator3/qr3_95/test.err.log
